# etc/cron.d/camviewers
* * * * * /etc/init.d/camviewers repair
