# etc/cron.d/camviewers
* * * * * for i in 1 2 3 4 5 6; do sleep 10; /etc/init.d/camviewers repair; done
